# hawkeye.crontab
#15 5 * * *  /home/bshaver/Documents/scripts/backup-home.sh > /tmp/backup-home.log 2>&1 
42 5 * * 6 /home/bshaver/machine-update.pl -m
